# deploy/rvsf-supervisord.conf
[unix_http_server]
file=/workspace/rvsf-supervisor.sock
chmod=0700

[supervisord]
logfile=/workspace/rvsf-supervisord.log
pidfile=/workspace/rvsf-supervisord.pid
childlogdir=/workspace

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///workspace/rvsf-supervisor.sock

[program:rvsf-api]
command=/bin/bash /workspace/radeon-voice-skill-foundry-current/scripts/radeon_api_supervisor.sh
directory=/workspace/radeon-voice-skill-foundry-current
autostart=true
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/workspace/rvsf-live-api.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:rvsf-tunnel]
command=/bin/bash /workspace/radeon-voice-skill-foundry-current/scripts/radeon_tunnel_supervisor.sh
directory=/workspace/radeon-voice-skill-foundry-current
autostart=true
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/workspace/rvsf-cloudflared.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:rvsf-origin-registrar]
command=/bin/bash /workspace/radeon-voice-skill-foundry-current/scripts/radeon_origin_registrar.sh watch
directory=/workspace/radeon-voice-skill-foundry-current
autostart=true
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/workspace/rvsf-origin-registrar.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
